# hidden-cam-detector.service
[Unit]
Description=Hidden Camera Detector (Pi Zero 2 W)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/hidden-cam-detector
EnvironmentFile=/home/pi/hidden-cam-detector/.env
ExecStart=/usr/bin/python3 /home/pi/hidden-cam-detector/main_system.py
Restart=always
RestartSec=2
KillSignal=SIGINT
TimeoutStopSec=8
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
